# deploy/systemd/nexus-api.service
[Unit]
Description=Nexus API Backend
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/surgery/Workspace/Nexus
Environment=PYTHONUNBUFFERED=1
Environment=HTTP_PROXY=http://127.0.0.1:10809
Environment=HTTPS_PROXY=http://127.0.0.1:10809
Environment=NO_PROXY=127.0.0.1,localhost,::1,100.121.67.94,100.87.152.11
ExecStart=/home/surgery/miniconda3/envs/ai_assist/bin/python -m nexus serve --host 100.121.67.94 --port 18100
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
